# license-server.service
[Unit]
Description=FLVX License Server
After=network.target

[Service]
Type=simple
Environment="SERVER_ADDR=:51888"
Environment="DB_PATH=/var/lib/license-server/license.db"
Environment="JWT_SECRET=change-me"
Environment="ADMIN_USERNAME=admin"
Environment="ADMIN_PASSWORD=change-me"
ExecStart=/usr/local/bin/license-server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
